# server/serverbot.service
[Unit]
Description=Telegram Pibot Server

[Service]
ExecStart=/srv/pibot/server/serverbot

[Install]
WantedBy=multi-user.target